# roles/common/files/etc/update-motd.d/95-ansible
#!/bin/sh

echo "___________                    ____  __.__  __         .__                 "
echo "\_   _____/___  __ _________  |    |/ _|__|/  |_  ____ |  |__   ____   ____   ______"
echo " |    __)/  _ \|  |  \_  __ \ |      < |  \   __\/ ___\|  |  \_/ __ \ /    \ /  ___/"
echo " |     \(  <_> )  |  /|  | \/ |    |  \|  ||  | \  \___|   Y  \  ___/|   |  \\___ \ "
echo " \___  / \____/|____/ |__|    |____|__ \__||__|  \___  >___|  /\___  >___|  /____  >"
echo "     \/                               \/             \/     \/     \/     \/     \/"
echo "________                  _________                                                 "
echo "\______ \   _______  __  /   _____/ ______________  __ ___________                  "
echo " |    |  \_/ __ \  \/ /  \_____  \_/ __ \_  __ \  \/ // __ \_  __ \                 "
echo " |    \`   \  ___/\   /   /        \  ___/|  | \/\   /\  ___/|  | \/                 "
echo "/_______  /\___  >\_/   /_______  /\___  >__|    \_/  \___  >__|                    "
echo "        \/     \/               \/     \/                 \/"
