# roles/common/files/etc/update-motd.d/95-ansible
#!/bin/sh
echo "\033[32m"

echo "$$$$$$$$\  $$$$$$\  $$$$$$$\   $$$$$$\"
echo "\__$$  __|$$  __$$\ $$  __$$\ $$  __$$\"
echo "   $$ |   $$ /  \__|$$ |  $$ |$$ /  \__|"
echo "   $$ |   $$ |      $$ |  $$ |$$ |"
echo "   $$ |   $$ |      $$ |  $$ |$$ |"
echo "   $$ |   $$ |  $$\ $$ |  $$ |$$ |  $$\"
echo "   $$ |   \$$$$$$  |$$$$$$$  |\$$$$$$  |"
echo "   \__|    \______/ \_______/  \______/"
echo ""
echo ""
echo ""
echo "$$\   $$\                           $$\ $$\"
echo "$$ |  $$ |                          $$ |$$ |"
echo "$$ |  $$ | $$$$$$\   $$$$$$\   $$$$$$$ |$$ | $$$$$$\   $$$$$$$\  $$$$$$$\"
echo "$$$$$$$$ |$$  __$$\  \____$$\ $$  __$$ |$$ |$$  __$$\ $$  _____|$$  _____|"
echo "$$  __$$ |$$$$$$$$ | $$$$$$$ |$$ /  $$ |$$ |$$$$$$$$ |\$$$$$$\  \$$$$$$\"
echo "$$ |  $$ |$$   ____|$$  __$$ |$$ |  $$ |$$ |$$   ____| \____$$\  \____$$\"
echo "$$ |  $$ |\$$$$$$$\ \$$$$$$$ |\$$$$$$$ |$$ |\$$$$$$$\ $$$$$$$  |$$$$$$$  |"
echo "\__|  \__| \_______| \_______| \_______|\__| \_______|\_______/ \_______/"
echo ""
echo ""
echo ""
echo "$$$$$$$$\                 $$\           $$\"
echo "\__$$  __|                \__|          \__|"
echo "   $$ | $$$$$$\  $$$$$$\  $$\ $$$$$$$\  $$\ $$$$$$$\   $$$$$$\"
echo "   $$ |$$  __$$\ \____$$\ $$ |$$  __$$\ $$ |$$  __$$\ $$  __$$\"
echo "   $$ |$$ |  \__|$$$$$$$ |$$ |$$ |  $$ |$$ |$$ |  $$ |$$ /  $$ |"
echo "   $$ |$$ |     $$  __$$ |$$ |$$ |  $$ |$$ |$$ |  $$ |$$ |  $$ |"
echo "   $$ |$$ |     \$$$$$$$ |$$ |$$ |  $$ |$$ |$$ |  $$ |\$$$$$$$ |"
echo "   \__|\__|      \_______|\__|\__|  \__|\__|\__|  \__| \____$$ |"
echo "                                                      $$\   $$ |"
echo "                                                      \$$$$$$  |"
echo "                                                       \______/"

echo "\033[m"
