# roles/common/files/etc/update-motd.d/95-ansible
#!/bin/sh
echo "\033[32m"

echo "db   db d88888b  .d8b.  d8888b. "
echo "88   88 88'     d8' \`8b 88  \`8D "
echo "88ooo88 88ooooo 88ooo88 88   88 "
echo "88~~~88 88~~~~~ 88~~~88 88   88 "
echo "88   88 88.     88   88 88  .8D "
echo "YP   YP Y88888P YP   YP Y8888D' "
echo "                                "
echo "db      d88888b .d8888. .d8888. "
echo "88      88'     88'  YP 88'  YP "
echo "88      88ooooo \`8bo.   \`8bo.   "
echo "88      88~~~~~   \`Y8b.   \`Y8b. "
echo "88booo. 88.     db   8D db   8D "
echo "Y88888P Y88888P \`8888Y' \`8888Y' "
echo "                                "
echo "db       .d8b.                  "
echo "88      d8' \`8b                 "
echo "88      88ooo88                 "
echo "88      88~~~88                 "
echo "88booo. 88   88                 "
echo "Y88888P YP   YP                 "
echo "\033[m"

